# bin/nagios-business-process-addon-0.9.6/etc/cron.d/nagiosbp.in
*/10 * * * * @apache_user@ @bindir@/nagios_bp_session_timeout 2>&1 | logger -t nagios_bp_session_timeout
